# examples/opensensor-space/systemd/opensensor_space_systemd.service
[Unit]
Description=OpenSensor Space Data Collection
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/enviroplus-python
Environment="PATH=/home/pi/enviroplus-python/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Main data collection service
ExecStart=/home/pi/enviroplus-python/.venv/bin/python examples/opensensor_space.py

# Restart policy
Restart=always
RestartSec=60

# Logging
StandardOutput=append:/home/pi/enviroplus-python/logs/opensensor_space.log
StandardError=append:/home/pi/enviroplus-python/logs/opensensor_space_error.log

# Security hardening (optional - uncomment if needed)
# NoNewPrivileges=true
# PrivateTmp=true

[Install]
WantedBy=multi-user.target
